# templates/service/neo4j.service.erb
[Unit]
Description=Neo4j Graph Database
After=network.target

[Service]
ExecStart=<%= scope['neo4j_home'] %>/bin/neo4j console
Restart=on-failure
User=<%= scope['user'] %>
Environment="NEO4J_CONF=<%= scope['config_dir'] %>" "NEO4J_HOME=<%= scope['data_prefix'] %>"
LimitNOFILE=<%= scope['service_ulimit'] %>
TimeoutSec=<%= scope['service_shutdown_timeout'] %>

[Install]
WantedBy=multi-user.target
